# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=npm start

[program:client]
command=/usr/local/lib/node_modules/http-server/bin/http-server ./